# stiefelsystem/etc/stiefelsystem.service
[Unit]
Description=Serve disks via Stiefelsystem

[Service]
Type=simple
ExecStart=/usr/bin/stiefelctl server

[Install]
WantedBy=multi-user.target
