# stiefelsystem/platform/files/initrd/etc/systemd/system/console-getty.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -f root -- \\u' --skip-login --noclear --keep-baud console 115200,38400,9600 $TERM
