# stiefelsystem/platform/files/initrd/etc/systemd/system/stiefel-server.service
[Unit]
Description=stiefelsystem server
Requires=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/stiefel-server
StandardOutput=journal+console
